# image/nvidia-cdi.service
[Unit]
Description=Generate NVIDIA CDI spec for nested Docker (when GPUs are present)
# The lab container receives GPUs via the host's CDI injection, so /dev/nvidiactl exists only when a
# GPU is attached. Generate an INNER CDI spec before docker starts so nested containers can request
# the GPU with `--device nvidia.com/gpu=all` (or `--gpus all`). No-op on CPU-only labs.
Before=docker.service
After=local-fs.target
ConditionPathExists=/dev/nvidiactl

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/mkdir -p /etc/cdi
# `|| true` so a transient generation failure never blocks docker/boot on a CPU-only or oddly-wired node.
ExecStart=/bin/sh -c '/usr/bin/nvidia-ctk cdi generate --output=/etc/cdi/nvidia.yaml || true'

[Install]
WantedBy=multi-user.target
