# image/ssh-hostkeys.service
[Unit]
Description=Generate per-container SSH host keys on first boot
# Run before sshd, and only if no host keys exist yet (the image ships without them so every
# container gets its own unique keys rather than sharing baked-in ones).
Before=ssh.service
ConditionPathExists=!/etc/ssh/ssh_host_ed25519_key

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/ssh-keygen -A

[Install]
WantedBy=multi-user.target
